# services/hodorbot/run
#!/bin/bash
umask 000

cd /opt/hodorbot
exec /sbin/setuser hodor ./rtmbot.py
